# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:mongodb]
command=/scripts/mongodb.sh
priority=1
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:redis]
command=/scripts/redis.sh
priority=1
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:freshclam]
command=/scripts/freshclam.sh
priority=2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:clamd]
command=/scripts/clamd.sh
priority=3
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

;[program:spamd]
;command=/scripts/spamd.sh
;priority=2
;autostart=false
;autorestart=true
;redirect_stderr=true
;stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:amavis]
command=/scripts/amavis.sh
priority=5
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:postgrey]
command=/scripts/postgrey.sh
priority=5
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:postfix]
command=/scripts/postfix.sh
priority=10
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:rs-admin]
command=/scripts/rs-admin.sh
priority=20
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

[program:start]
command=/scripts/start-services.sh
priority=999
startretries=1
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log

